# root/app/donation/supervisord.conf
[program:worker]
command = /app/donation/worker.sh
autorestart = true
pidfile = /run/worker.pid
priority = 2
